# scripts/go-eva.service
[Unit]
Description=go-eva Audio Spatial Awareness Daemon
Documentation=https://github.com/teslashibe/go-eva
After=network.target reachy-mini-daemon.service
Wants=reachy-mini-daemon.service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/go-eva -config /etc/go-eva/config.yaml
Restart=always
RestartSec=5
WatchdogSec=30

# Logging to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=go-eva

# Give time for USB device to initialize
ExecStartPre=/bin/sleep 2

# Environment
Environment=HOME=/root

# Resource limits
MemoryMax=64M
CPUQuota=15%

[Install]
WantedBy=multi-user.target
